# deployment/systemd/aura-daemon.service
[Unit]
Description=AURA daemon (user telemetry producer)
Documentation=https://github.com/zapsaang/aura
After=default-target

[Service]
Type=simple
Environment=AURA_HEARTBEAT_MS=500
Environment=AURA_SHM_PATH=%t/aura_state.dat
Environment=RUST_LOG=info
ExecStart=%h/.local/bin/aura-daemon --heartbeat-ms ${AURA_HEARTBEAT_MS} --shm-path ${AURA_SHM_PATH}
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
